# bin/webshell.conf
# Upstart script
# /etc/init/webshell.conf

description "Web TTY"
author      "Webshell"

start on started mountall
stop on shutdown

respawn
respawn limit 20 5

exec sudo -u root webshell -p 3000
