# bin/webshell.service
# systemd unit file
#
# place in /etc/systemd/system
# systemctl enable webshell.service
# systemctl start webshell.service

[Unit]
Description=Webshell Daemon
After=network.target

[Service]
User=root
Group=root

WorkingDirectory=/opt/webshell
ExecStart=webshell-server -p 9443

[Install]
WantedBy=multi-user.target
